# qubes/tb-dvm/rw/config/rc.local
#!/usr/bin/env bash

set -euo pipefail

/rw/config/ram-drive-setup
/rw/config/tb-customize
