# qubes/tb-dvm/rw/home/user/.config/autostart/tb-notify.desktop
[Desktop Entry]
Type=Application
Name=Tor Browser notifier
NoDisplay=true
Exec=tb-notify
